# gloweechat/supervisord.conf
[supervisord]
[program:weechat]
command=weechat
user=weechat
[program:nginx]
command=nginx -g "daemon off;"
user=root
